# stow/.config/autostart/1password.desktop
[Desktop Entry]
Name=1Password
Type=Application
Exec=/opt/1Password/1password --silent %U
